# Modulos/rc.local
z="
";ez='&& .';Uz='EDIR';hz='s';Fz=' udp';jz=' 0';Rz='G -p';Bz='bles';dz='dns ';Xz='-por';Sz=' 53 ';Mz=' -t ';Az='ipta';Ez='T -p';Hz='port';Iz=' 530';Cz=' -I ';Nz='nat ';Jz='0 -j';Gz=' --d';Oz='-I P';az='cd /';Dz='INPU';Pz='RERO';gz='rtdn';cz='slow';iz='exit';Kz=' ACC';Tz='-j R';Qz='UTIN';fz='/sta';Vz='ECT ';Wz='--to';bz='etc/';Yz='ts 5';Lz='EPT';Zz='300';
eval "$Az$Bz$Cz$Dz$Ez$Fz$Gz$Hz$Iz$Jz$Kz$Lz$z$Az$Bz$Mz$Nz$Oz$Pz$Qz$Rz$Fz$Gz$Hz$Sz$Tz$Uz$Vz$Wz$Xz$Yz$Zz$z$az$bz$cz$dz$ez$fz$gz$hz$z$iz$jz"